# samus/etc/systemd/system/pacman-rankmirrors.timer
[Unit]
Description=Run pacman rankmirrors weekly

[Timer]
OnCalendar=weekly
RandomizedDelaySec=12h
Persistent=true

[Install]
WantedBy=timers.target
