# moln/systemd/snapshot-weekly.service
[Unit]
Description=Weekly Snapshots
Wants=timer-weekly.timer

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/zfs-snapshot.rb -s data -k 4 -p weekly -S "%%s"

[Install]
WantedBy=timer-weekly.target
